# scripts/deployment/75flow.service
[Unit]
Description=75Flow Django Application
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/project
Environment="PATH=/path/to/your/project/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=flow75.settings_production"
Environment="ENVIRONMENT=production"
Environment="SECRET_KEY=your-secret-key-here"
Environment="ALLOWED_HOSTS=shrubuilds.in,www.shrubuilds.in"
Environment="CSRF_TRUSTED_ORIGINS=https://shrubuilds.in,https://www.shrubuilds.in"
ExecStart=/path/to/your/project/venv/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 --access-logfile /var/log/75flow/access.log --error-logfile /var/log/75flow/error.log flow75.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 